# config/supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0777

[supervisord]
logfile=/dev/null
pidfile=/var/run/supervisord.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:squid]
command=/usr/sbin/squid -N -d -f /config/squid/squid.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/squid.out.log
stderr_logfile=/var/log/squid.err.log

